# deploy/usstock-frontend.service
[Unit]
Description=US Stock Trading Frontend (static serve)
After=network.target usstock-backend.service

[Service]
Type=simple
User=chans
WorkingDirectory=/home/chans/us-stock/frontend
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npx serve -s dist -l 3001
Restart=always
RestartSec=10
Environment=HOME=/home/chans
MemoryMax=384M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
